# com.finhire-korea.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finhire-korea.daily-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/rudolph/Documents/Codex/2026-07-22/referenced-chatgpt-conversation-this-is-untrusted/sync_daily.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/rudolph/Documents/Codex/2026-07-22/referenced-chatgpt-conversation-this-is-untrusted</string>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key>
  <string>/Users/rudolph/Documents/Codex/2026-07-22/referenced-chatgpt-conversation-this-is-untrusted/finhire-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rudolph/Documents/Codex/2026-07-22/referenced-chatgpt-conversation-this-is-untrusted/finhire-sync-error.log</string>
</dict>
</plist>
